# sniffer-agent.ini
[program:mysql-sniffer-agent]
directory=/app/
command=/app/sniffer-agent --interface=eth0 --port=30006 --log_level=debug --export_type=api --api-report-url=http://172.24.106.61:9999/v1/mysql/instances/exec-sql
autostart=true
autorestart=true
startretries=3
stderr_logfile=/tmp/test_stderr.log
stdout_logfile=/tmp/test_stdout.log
user=root